# usr/lib/systemd/user/gamescope-session.target
[Unit]
Description=Gamescope Session Target